# pkg/pkgbuild/doubao-voice-input-electron.service
[Unit]
Description=Doubao Voice Input - Real-time speech-to-text
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=/opt/doubao-voice-input-electron/doubao-voice-input-electron

[Install]
WantedBy=graphical-session.target
